# public_projects/ai-usage-dashboard/install/usage-dashboard.service
# AI Usage Dashboard systemd service
# Install via: sudo bash install/install.sh /absolute/path/to/project
# Or manually: edit PROJECT_DIR and USER below, then:
#   sudo cp install/usage-dashboard.service /etc/systemd/system/usage-dashboard.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now usage-dashboard

[Unit]
Description=AI Usage Dashboard
Documentation=https://github.com/anomalyco/opencode
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/usage-dashboard
ExecStart=PROJECT_DIR/venv/bin/python3 -m main
WorkingDirectory=PROJECT_DIR/backend
Restart=on-failure
RestartSec=5
User=PROJECT_USER
Group=PROJECT_USER
StandardOutput=journal
StandardError=journal

# Allow access to user home dirs (AGY DBs, OpenCode config, etc.)
ProtectHome=false
ProtectSystem=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
